# data/console-kit-log-system-start.service.in
[Unit]
Description=Console System Startup Logging
DefaultDependencies=no
Wants=sysinit.target
After=sysinit.target
Before=shutdown.service

[Service]
Type=oneshot
ExecStart=@sbindir@/ck-log-system-start
RemainAfterExit=yes
