# data/console-kit-log-system-stop.service.in
[Unit]
Description=Console System Shutdown Logging
DefaultDependencies=no
Wants=sysinit.target
After=sysinit.target console-kit-system-start.service
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=@sbindir@/ck-log-system-stop
