# data/console-kit-log-system-restart.service.in
[Unit]
Description=Console System Reboot Logging
DefaultDependencies=no
Wants=sysinit.target
After=sysinit.target console-kit-system-start.service
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=@sbindir@/ck-log-system-restart
