# data/console-kit-log-system-stop.service.in
[Unit]
Description=Console System Shutdown Logging
DefaultDependencies=no
Wants=local-fs.target sysinit.target
After=local-fs.target sysinit.target console-kit-system-start.service
Before=poweroff.service reboot.service halt.service killall.service

[Service]
Type=oneshot
ExecStart=@sbindir@/ck-log-system-stop
