# data/console-kit-log-system-start.service.in
[Unit]
Description=Console System Startup Logging
DefaultDependencies=no
Wants=local-fs.target sysinit.target
After=local-fs.target sysinit.target
Before=poweroff.service reboot.service halt.service killall.service

[Service]
Type=oneshot
ExecStart=@sbindir@/ck-log-system-start
RemainAfterExit=yes
